# service/monitoring/run
#!/bin/bash

packages=(influxdb telegraf grafana zabbix)
load_package_run_file "${packages[*]}"