# service/monitoring/run
#!/bin/bash

packages=(influxdb telegraf grafana)
load_package_run_file "${packages[*]}"