# service/dns/run
#!/bin/bash

packages=(bind)
load_package_run_file "${packages[*]}"